# simple_http_server.conf
[program:simple_http_server]
directory=/simple_http_server
command=python3 main.py
autostart=true
autorestart=true